# packages/wolkenbruch/wolkenbruch-0.6.3.tar.gz/wolkenbruch-0.6.3/extra/systemd/wolkenbruch.service
[Unit]
Description = Check whether it’s likely to rain and send an e-mail reminder

[Service]
ExecStart = /usr/bin/wolkenbruch
