# packages/odoopbx/odoopbx-0.88.tar.gz/odoopbx-0.88/salt/roots/agent/agent.service
[Unit]
Description=OdooPBX Agent
After=network.target

[Service]
Environment=LANG={{ agent.get('locale', 'C.UTF-8') }}
{% if salt['environ.get']('VIRTUAL_ENV') %}
Environment="PATH={{ grains['systempath']|join(':')}}:{{ salt['environ.get']('VIRTUAL_ENV') }}/bin"
{% else %}
Environment="PATH={{ grains['systempath']|join(':')}}"
{% endif %}
Type=simple
User=root
ExecStart=/usr/bin/env odoopbx run agent
Restart=on-failure
RestartSec=10
RemainAfterExit=no
StandardOutput=null
StandardError=null
SyslogIdentifier=odoopbx-agent

[Install]
WantedBy=multi-user.target
