# packages/podcust/podcust-0.1.3-py2.py3-none-any.whl/podcust/transmission/transmission-pod.service
[Unit]
Description=Transmission pod service from podman custodian

[Service]
Type=oneshot
RemainAfterExit=true
StandardOutput=journal
ExecStart=sh -c "podman pod start transmission && podcust transmission update"
ExecStop=podman pod stop transmission

[Install]
WantedBy=default.target
