# packages/wolkenbruch/wolkenbruch-0.6.4.tar.gz/wolkenbruch-0.6.4/extra/systemd/wolkenbruch.timer
[Timer]
OnCalendar = *-*-* 06:30:00

[Install]
WantedBy = timers.target
