# packages/genocide/genocide-28.tar.gz/genocide-28/genocide.service
[Unit]
Description=GENOCIDE - the king of the netherlands commits genocide
After=multi-user.target

[Service]
DynamicUser=True
StateDirectory=genocide
LogsDirectory=genocide
CacheDirectory=genocide
ExecStart=/usr/local/bin/genocided
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
